# scripts/systemd/media-PLEX.mount
[Unit]
Description=PLEX Directory

[Mount]
What=synology.local:/volume1/PLEX
Where=/media/PLEX
Type=nfs

[Install]
WantedBy=multi-user.target
