# scripts/systemd/media-home.automount
[Unit]
Description=Automount Home Directory

[Automount]
Where=/media/home

[Install]
WantedBy=multi-user.target
