# scripts/systemd/media-PLEX.automount
[Unit]
Description=Automount PLEX Directory

[Automount]
Where=/media/PLEX

[Install]
WantedBy=multi-user.target
